# config/systemd/user/app-aw\x2dqt@autostart.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/aw-qt --no-gui
