# config/systemd/user/backup-package-list.timer
[Unit]
Description=Back up package list every 12h

[Timer]
OnBootSec=5min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
