# config/systemd/user/backup-passwords.service
[Unit]
Description=Back up passwords to remote repo

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'pass git pull && pass git push'
