# config/systemd/user/background.timer
[Unit]
Description=Change Unsplash wallpaper twice a day

[Timer]
OnCalendar=00,12:00
Persistent=true

[Install]
WantedBy=timers.target
