# config/systemd/user/background.timer
[Unit]
Description=Change Unsplash wallpaper twice a day

[Timer]
OnStartupSec=1h
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
