# config/systemd/user/backup-passwords.timer
[Unit]
Description=Back up passwords every 12h

[Timer]
OnCalendar=00,12:00
Persistent=true

[Install]
WantedBy=timers.target
