# config/systemd/user/backup-passwords.timer
[Unit]
Description=Back up passwords every 12h

[Timer]
OnBootSec=5min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
