# config/systemd/user/backup-package-list.timer
[Unit]
Description=Back up package list every 12h

[Timer]
OnCalendar=00,12:00
Persistent=true

[Install]
WantedBy=timers.target
