# config/systemd/user/reminder.service
[Unit]
Description=Reminder of calendar events

[Service]
Type=oneshot
ExecStart=%h/bin/calendar-remind
