# config/systemd/user/tad.service
[Unit]
Description=Time Awareness Daemon

[Service]
Type=simple
ExecStart=%h/bin/tad

[Install]
WantedBy=default.target
