# config/systemd/user/tad.service
[Unit]
Description=Time Awareness Daemon

[Service]
Type=simple
ExecStart=/home/procrat/bin/tad

[Install]
WantedBy=default.target
