# config/systemd/user/watchdog.service
[Unit]
Description=Watchdog for system failures

[Service]
Type=oneshot
ExecStart=%h/bin/systemd-watchdog
